# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=critical
pidfile=/run/supervisord.pid
user=root

[unix_http_server]
file=/run/supervisor.sock
username=supervisor
password=%(ENV_SUPERVISOR_PASSWORD)s

[supervisorctl]
serverurl=unix:///run/supervisor.sock
username=supervisor
password=%(ENV_SUPERVISOR_PASSWORD)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:unbound]
command=/usr/sbin/unbound -d -c /etc/unbound/unbound.conf
autostart=true
autorestart=true
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
priority=0

[program:redis]
command=/usr/bin/redis-server --daemonize no --dir /app/data/redis --save 60 1
autostart=true
autorestart=true
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
priority=1

[program:rspamd]
command=/usr/bin/rspamd -f -u _rspamd -g _rspamd
autostart=true
autorestart=true
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
priority=2

[program:boxwatchr]
command=python3 /app/main.py
user=%(ENV_PUID)s
autostart=true
autorestart=unexpected
exitcodes=0,2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3